# web/rootfs/etc/supervisor.d/custom.ini
[supervisord]
nodaemon=true

[program:php-fpm]
command = /usr/sbin/php-fpm7
user = root
stdout_events_enabled=true
stderr_events_enabled=true
autostart = true

[program:nginx]
command = /usr/sbin/nginx
user = root
stdout_events_enabled=true
stderr_events_enabled=true
autostart = true

